# confs/supervisord/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=www-data

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.0 -c /etc/php7.0/fpm
stdout_logfile=/srv/log/php-fpm-stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/srv/log/php-fpm-stderr.log
stderr_logfile_maxbytes=0
exitcodes=0
autostart=true
autorestart=true
priority=5
user=www-data

[program:pdns]
stdout_logfile=/srv/log/pdns-stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/srv/log/pdns-stderr.log
stderr_logfile_maxbytes=0
exitcodes=0
autostart=true
autorestart=true
priority=5
user=www-data
command=/usr/sbin/pdns_server --guardian=yes --socket-dir=/srv/run/pdns --setgid=www-data --setuid=www-data --launch=random --loglevel=9
